# deploy/supervisor.conf
; Phase 12 — supervisor cho deploy không-Docker. Copy vào /etc/supervisor/conf.d/
; sudo supervisorctl reread && sudo supervisorctl update

[program:taskboard-queue]
command=php /var/www/taskboard/laravel-backend/artisan queue:work --tries=3 --max-time=3600
autostart=true
autorestart=true
user=www-data
numprocs=2
process_name=%(program_name)s_%(process_num)02d
stopwaitsecs=3600
redirect_stderr=true
stdout_logfile=/var/log/taskboard/queue.log

[program:taskboard-reverb]
command=php /var/www/taskboard/laravel-backend/artisan reverb:start --host=0.0.0.0 --port=8080
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/taskboard/reverb.log

[program:taskboard-scheduler]
command=sh -c "while true; do php /var/www/taskboard/laravel-backend/artisan schedule:run --no-interaction; sleep 60; done"
autostart=true
autorestart=true
user=www-data
redirect_stderr=true
stdout_logfile=/var/log/taskboard/scheduler.log
